# templates/initscript.erb
#!/bin/bash
# chkconfig: 345 20 99
# description: <%= @description %>
#. /etc/rc.d/init.d/functions
# Default-Start: 3 5
# Default-Stop: 0 6

<% if @option_scripts.any? -%>
<% @option_scripts.each do |val| -%>
. <%= val %> > /dev/null 2>&1
<% end -%>
<% end -%>

PIDFILE="<%= @pid_file %>"
<% if defined?(@cmd) -%>
CMD="<%= @cmd %>"
<% end -%>

<% if @run_user == "root" -%>
EXECSTART=""
<% else -%>
RUN_USER=<%= @run_user %>
if [ "$(id -u)" -eq "$(id -u <%= @run_user %>)" ];
then
  EXECSTART="/bin/bash"
else
  EXECSTART="/bin/su -p -s /bin/bash $RUN_USER"
fi
<% end -%>

case $1 in
  'start')
    <%- if defined?(@cmd) -%>
    $EXECSTART $CMD &
    echo $! > $PIDFILE
    <%- else -%>
    <%= @cmd_start %>
    <%- end -%>
    exit 0
  ;;
  'stop')
  <%- if defined?(@cmd) -%>
  if [ -z "$(cat $PIDFILE)" ];
  then
    echo "stopped"
    exit 1
  else
    kill $(cat $PIDFILE)
    sleep 1
    kill -9 $(cat $PIDFILE)
    sleep 1
    rm $PIDFILE
    exit 0
  fi
  <%- else -%>
  $EXECSTART <%= @cmd_stop %>
  <%- end -%>
  ;;
  'restart')
    $0 stop
    $0 start
  ;;
  'status')
  <%- if defined?(@tcp_listen_check) -%>
  LSOFBIN=$(which lsof 2> /dev/null)

  if [ -z "${LSOFBIN}" ];
  then
    NETSTATBIN=$(which netstat 2> /dev/null)

    if [ -z "${NETSTATBIN}" ];
    then
      echo "neither lsof nor netstat found, aborting"
      exit 1
    fi

    netstat -tpln | grep ":<%= @tcp_listen_check %>\b" > /dev/null 2>&1

    if [ $? -eq 0 ];
    then
      echo "running"
      exit 0
    else
      echo "stopped"
      exit 1
    fi
  fi

  $LSOFBIN -i :<%= @tcp_listen_check %> -sTCP:LISTEN > /dev/null 2>&1

  if [ $? -eq 0 ];
  then
    echo "running"
    exit 0
  else
    echo "stopped"
    exit 1
  fi
  <%- else -%>
    if [ -z "$(cat $PIDFILE 2>/dev/null)" ];
    then
      echo "stopped"
      exit 1
    else
      if [ $(ps -p $(cat $PIDFILE) | wc -l 2>/dev/null) -eq 2 ];
      then
        echo "running"
        exit 0
      else
        echo "stopped"
        exit 1
      fi
    fi
  <%- end -%>
  ;;
  *)
  echo "$0 $1: unrecognized command"
  exit 1
  ;;
esac
